# output/build/fontconfig-2.12.1/.br_filelist_before
011269dec3d3814d3632065d224dc758  ./lib/libc-2.23.so
019878bc2aea3526dc97d09ca29c40d2  ./usr/bin/chocolate-heretic
0297048659feb06aac253a1287c252f4  ./usr/share/alsa/ucm/SDP4430/SDP4430.conf
035ab192ff65e43c0d2f9afd94484117  ./usr/share/alsa/cards/VX222.conf
03a3f214c7192dd785a767e2bb2e7b8a  ./lib/libm-2.23.so
05a1f97b4d9334371aea13710d72d024  ./usr/share/alsa/cards/aliases.alisp
05a964fbefa0055acb78ccbd9b73fbed  ./usr/lib/libswresample.so.2.7.100
0691559eddbbb03c8e81f5cbda329f27  ./usr/bin/chocolate-server
090dde7adf8e77be248f27e3ae615e46  ./etc/group
0a39be7f9b4af4c9f71926a7dd978362  ./usr/share/ffmpeg/libvpx-1080p.ffpreset
0b375c3957371c0ac5f2bbcb5100c845  ./usr/share/alsa/sndo-mixer.alisp
0ef862df711343b6428d9b1e5692bb39  ./usr/bin/chocolate-strife
0f6f28e5fc5e27299332f4d59d3ef1bc  ./usr/share/alsa/cards/aliases.conf
0fed7b03412669f1daaf03cd15260e0c  ./usr/share/alsa/ucm/PandaBoardES/hifi
0ff98410642ae96284cedfe843df8079  ./usr/share/alsa/pcm/side.conf
11579a5b50f3b370e902afd460cb149a  ./usr/share/alsa/pcm/iec958.conf
11eeee86df9bfb1edcd2f5ac558ed871  ./usr/share/alsa/init/hda
120a28e52423409a5399346b4b0f937c  ./usr/include/expat_config.h
125b886122712510e527bbd2e211a5a0  ./usr/share/alsa/ucm/PandaBoardES/hifiLP
132f6976651aa6f8b56e73d015bf4da2  ./usr/share/alsa/cards/ICH4.conf
13ce51bab20504473d858ccb7fb7c845  ./usr/share/alsa/ucm/skylake-rt286/skylake-rt286.conf
146dd6077999cf9557891c91fa27d36a  ./lib/librt-2.23.so
15777f2c684ae84c0024ba4bcadb7ddb  ./usr/share/alsa/cards/FM801.conf
1837782a0132e61065179c8ceb1172b6  ./usr/share/alsa/cards/ICH-MODEM.conf
189aaa6830dd210b9389d552fb9a807a  ./lib/libatomic.so.1.2.0
18b5c35ccedb88528eea3c8a7e920f5c  ./usr/lib/libswscale.so.4.6.100
19001a1f98e733cb03da6908b8b831e6  ./usr/share/alsa/cards/ICE1712.conf
19a5f2f42aa2a86eca1365c491877ea0  ./usr/share/alsa/ucm/PandaBoard/FMAnalog
1a78fbc2f9a7f33e83676f30dd3990e3  ./usr/share/alsa/cards/GUS.conf
1a810ab11d9ab76f229b071ef575695f  ./usr/lib/libSDL_net-1.2.so.0.8.0
1a8aa0603e08dfe8a464918eddf9b629  ./usr/share/alsa/pcm/surround51.conf
1aefd1542930af5f60dfd687685896d8  ./usr/share/alsa/init/help
1af0d2b3b63e57d2a32ed6f1f7203ad1  ./usr/share/alsa/cards/CMI8738-MC8.conf
1c08ebdc8396f3bc30f7b42cd5862064  ./usr/share/alsa/ucm/SDP4430/voiceCall
1cad5125d81c1f01ea675f0c31926502  ./usr/include/expat_external.h
1d95181ffb0d570a8eea7bb011752046  ./usr/lib/libmenu.so.6.0
1da78eba24dcb2710e6d93f3bccd84e1  ./usr/lib/libnl-nf-3.so.200.25.0
1e43e61b0e4fba70d63b5b6d3dce2a7c  ./usr/share/terminfo/a/ansi
218e0321b6a2cb27864d18834380f201  ./usr/share/alsa/cards/AU8810.conf
218e1b705fa98a3bf0b7a9048073727e  ./usr/share/alsa/cards/Loopback.conf
229bd6a21777499d446c359beea2b6e2  ./usr/share/alsa/ucm/PandaBoardES/voiceCall
23f525ff379e448abb3b133c63f9dff5  ./usr/share/alsa/ucm/VEYRON-I2S/HiFi.conf
25bdf86baca90aa75cd41ddfce94579b  ./usr/share/alsa/cards/YMF744.conf
26148ea963d68ba67ead72ddd61f383e  ./usr/share/alsa/topology/broadwell/broadwell.conf
26eecc0de2eba76bd2d31d0c02e14078  ./usr/share/alsa/ucm/tegraalc5632/tegraalc5632.conf
26ef45a8c3cb6d21f96f76d7c21afc5a  ./usr/share/alsa/topology/bxtrt298/bxt_i2s.conf
279b9101aa6c4d8f3e89d8be8537b8ec  ./usr/share/alsa/pcm/surround41.conf
285935c2ac102ac093265bacefa8c2e0  ./usr/bin/chocolate-doom-setup
285935c2ac102ac093265bacefa8c2e0  ./usr/bin/chocolate-heretic-setup
285935c2ac102ac093265bacefa8c2e0  ./usr/bin/chocolate-hexen-setup
285935c2ac102ac093265bacefa8c2e0  ./usr/bin/chocolate-strife-setup
297c6a4f51d6d3b526bb07bd478af00d  ./usr/share/alsa/cards/pistachio-card.conf
29f31d24345d79f75703e9f7e44d5b80  ./usr/share/alsa/cards/SI7018/sndop-mixer.alisp
29f68994c641fe4b1700d661eaf48d28  ./usr/share/alsa/cards/Audigy.conf
2ba852ea108988355e997fd7f485b15e  ./usr/share/alsa/init/test
2c2df1b9cb603f9c31c46162d6ac307f  ./usr/share/alsa/alsa.conf
2ddfb8c20a7ca55ff6c906911a857148  ./usr/share/terminfo/v/vt200
2ddfb8c20a7ca55ff6c906911a857148  ./usr/share/terminfo/v/vt220
2f55855f3b6319b2c3b3a087d01271c6  ./usr/share/alsa/ucm/DB410c/DB410c.conf
301f6407499f4b6fc14a42139dfce94f  ./usr/share/alsa/ucm/chtrt5645/HiFi.conf
31e978531fed0ab91bc96779d5dad3d6  ./usr/share/terminfo/l/linux
323ff2d327fcb45402a6d1d2a48f7072  ./usr/share/alsa/cards/NFORCE.conf
32ee939da56757aaee76e97f354f4762  ./usr/bin/ffmpeg
332ebea96715dc9ba66d66a13d0be3d6  ./usr/sbin/iw
33b7e08d7cd1ed759c9342a69094e674  ./usr/lib/libnl-genl-3.so.200.25.0
34a68ba0a48b8422e89264296fcea8e5  ./usr/share/alsa/ucm/broadwell-rt286/broadwell-rt286.conf
3505999adcc2d854772931dc2d33c247  ./usr/bin/dwebp
359ec096fae48c11cdb9adb8e113bc3c  ./etc/profile.d/umask.sh
36848341f4012291d662ee17937d9998  ./usr/bin/tinypcminfo
36a99b51f99af0d02d863895d73ca8d5  ./usr/share/ffmpeg/libvpx-720p.ffpreset
377a968489a2e822eab7ae9518d62097  ./usr/share/alsa/cards/SI7018/sndoc-mixer.alisp
3a7bcb7e3aca96c0f0e0ede436677e26  ./usr/share/alsa/cards/VIA8233A.conf
3b9fcbe59b633711d68173789a4a42d8  ./usr/share/alsa/cards/SI7018.conf
3ca597df8d194803821758ca48b9feea  ./usr/lib/libusb-1.0.so.0.1.0
3ce34797a1f01ebecdc86f4f68ce0237  ./usr/share/alsa/cards/CMI8738-MC6.conf
3df817d0e17aa2224a94c2d89bf8a88e  ./THIS_IS_NOT_YOUR_ROOT_FILESYSTEM
3e07259e58674631830b152e983ca995  ./etc/libnl/classid
3e1714fbe385fd21c95c75237766055d  ./usr/share/alsa/speaker-test/sample_map.csv
3ee6334fe2e6557eb85d748bbb23dd50  ./usr/share/alsa/ucm/PAZ00/HiFi.conf
3f3567b546158d9dfda52d69c45d606f  ./usr/share/alsa/cards/Aureon71.conf
40901c966f81ceeba77952cbed76d94c  ./usr/share/ffmpeg/libvpx-1080p50_60.ffpreset
40c3dac132e24e00fbfc51ccf961f72e  ./usr/lib/libfdt-1.4.4.so
40fef3764c096f653e60b1738707db2d  ./usr/bin/cwebp
438f9fd53d149a0046f7452e2b85bd3d  ./usr/sbin/alsactl
43e05287858104e07d062ed5a68d7293  ./usr/share/alsa/ucm/SDP4430/record
44fab56e772b1da4a9160a866d396833  ./usr/share/alsa/cards/PMacToonie.conf
458a10c34b46416c77745e47b8713a0d  ./usr/share/alsa/ucm/broadwell-rt286/HiFi
45f24b125e7cf866e648b4dafb0b255c  ./usr/lib/libSDL_image-1.2.so.0.8.4
471fe244e731c10bb3bbca1386bfd8b5  ./usr/share/alsa/ucm/chtrt5645/chtrt5645.conf
4760b042837ce9e23958037544b39831  ./usr/share/alsa/ucm/skylake-rt286/HiFi
487e82468bb038522279d2c1459a2120  ./usr/share/alsa/pcm/default.conf
48e72bd9ab4a4a69685f67344a1aa306  ./usr/lib/libasound.so.2.0.0
493a73269c19d8f3e498d612506c2303  ./usr/lib/libfreetype.so.6.14.0
4ec22ffcfa06e137629b0cb6a25ddb4c  ./usr/share/terminfo/x/xterm
4edcec0b89ff43eb3c2168fbcb66f007  ./usr/share/alsa/pcm/surround50.conf
4fcd0cc4b6d8c6f121fa9420ab7a95ab  ./usr/share/alsa/cards/PMac.conf
50d35b883737f91ed20806aca04eedec  ./usr/bin/tinyplay
518bc1c25bf80eb22d826ba5fd0bc3ca  ./usr/lib/libusb-0.1.so.4.4.4
51b785606f4beaf3f80b18f6d9e530ba  ./usr/lib/libturbojpeg.so.0.1.0
566bef0c9dbc473671208d82870f424f  ./usr/share/terminfo/p/putty
5670c78dab0d1d25f073720ea55fc481  ./bin/busybox
571ace2d72440b17ba2ea8d1568b6131  ./usr/share/alsa/cards/PS3.conf
5b82c82db0584ea54783c857562a5b12  ./lib/libpthread-2.23.so
5e9b1e683a121d59cf4065ccb16ddddf  ./usr/bin/tinycap
5fe25a4370c02f4ddee6285eac2ab974  ./usr/share/terminfo/v/vt102
6110ebdd56759c3c525eb7707009ed99  ./usr/share/alsa/ucm/DAISY-I2S/HiFi.conf
614f55c9c48f2a8cd4eb16a640e28f3c  ./usr/share/alsa/pcm/surround21.conf
61f1b3f5a9ddd7071295062ec45b0900  ./usr/share/alsa/cards/HdmiLpeAudio.conf
642da11e537ffd810ee5fe00d144b3eb  ./lib/libanl-2.23.so
648a1658a30f21b6972e59f511ced29e  ./usr/share/alsa/pcm/dmix.conf
65541d1b3ceee25cefac39f95c3b9b40  ./usr/share/alsa/cards/EMU10K1X.conf
65553a1d0273070536b3daf91618c372  ./etc/profile
65867601a12b7d422f6691587be46ca6  ./usr/share/alsa/ucm/PandaBoard/voiceCall
662558aa6933002d5d4e63018e293ad7  ./usr/share/alsa/cards/AU8830.conf
675b078c1d6e90cd6b95c71a58d47209  ./usr/share/alsa/ucm/SDP4430/voice
6805b001e1195fa59858107dcf8368cd  ./etc/inittab
68f280ddb7dcaeb23bcdb6ec344d6f39  ./usr/share/alsa/ucm/PandaBoard/hifi
6a8ee7e8c36df3885a327eb401cb5782  ./lib/libgcc_s.so.1
6b31f620c36bb2d6c54cfaf71cc5a086  ./etc/init.d/rcK
6c87c9ac4db746776a45520822e7eac9  ./usr/bin/mmc
6cb822c5baa87c0112813be59ad7ff67  ./usr/share/alsa/init/ca0106
6d24d92e31fa4d5a0fb8dce6c3ad6b12  ./etc/os-release
6e873ea02faa57e6997d61f4e4bb12ea  ./usr/lib/pkgconfig/expat.pc
6ea82e75711163776abde2d848a14afe  ./lib/libutil-2.23.so
6f0d8d5b0ed81c21725da1024a9cf0e9  ./usr/share/alsa/ucm/SDP4430/FMAnalog
6f33c79f46f584cee3171dc20fe898fe  ./usr/share/alsa/pcm/hdmi.conf
70021d405f76fafa35abab864dfe52f7  ./usr/lib/libncurses.so.6.0
711c6a4e5ec3f5e7088688d4aad01b75  ./etc/init.d/S20urandom
7277980cb323f7dcdcf4e89c7c0de379  ./usr/share/alsa/cards/VIA8237.conf
732b7677527e63ace921573027683274  ./etc/issue
735b1722b2d3bc7b3d5920966279f888  ./usr/share/alsa/cards/EMU10K1.conf
74e8d1be9fb6781f45f8a9740ae3f1e3  ./usr/share/terminfo/p/putty-vt100
751ca3865f2b0bc1467099f7919ec9b7  ./usr/share/alsa/ucm/DAISY-I2S/DAISY-I2S.conf
75ff0869b3f20382ebd5f629e0629bb8  ./usr/share/alsa/ucm/PandaBoard/PandaBoard.conf
7613dbc41b2dc3258195b6b6abd0f179  ./etc/libnl/pktloc
76f02a748149f1aa945aa418ea65b2bc  ./etc/init.d/rcS
78523c087c9a8c87cbf75df388fc6b60  ./usr/share/alsa/cards/ATIIXP-MODEM.conf
78669b30fab539d092021b534e07e8a6  ./usr/share/alsa/ucm/PAZ00/PAZ00.conf
79422c541750efcba29132115d4207a7  ./usr/share/alsa/ucm/PandaBoardES/PandaBoardES.conf
797d31702e99bd7c2767157905f0e0ed  ./usr/bin/amixer
799b745cf912a1fc3295781a227deeed  ./usr/bin/chocolate-hexen
7a7a3a47803b56c020709000fba94204  ./usr/lib/libSDL_mixer-1.2.so.0.12.0
7b51df8150ab9ee471a0d061c5a73a6c  ./usr/share/alsa/ucm/DB410c/HiFi
7b52aec959d8a0608f8dadbec6a6cbe3  ./usr/share/alsa/ucm/broxton-rt298/Hdmi1
7c04e662b37e09a98c46141cf08fb247  ./etc/passwd
7c200bd842745c4c789a1396eeebf6db  ./usr/share/alsa/cards/ENS1370.conf
7c708db5587cf53a6b4da28fb3118c01  ./usr/share/alsa/alsa.conf.d/README
7cec0ab7b91a337a0e585c1992ca3753  ./usr/lib/libexpat.so.1.6.4
7ec0863f1342719c134d2d7115fe6a6f  ./usr/share/alsa/ucm/SDP4430/hifiLP
80bfc68a72de195d830a79719ace8a04  ./usr/lib/libpanel.so.6.0
81cf06379072d0717eeea3dc019d2efa  ./usr/share/terminfo/x/xterm-xfree86
81e5f5cc52a658afd58700c84527fb15  ./usr/share/ffmpeg/ffprobe.xsd
84b05375889f2213644e3ac24820c5c4  ./usr/share/alsa/init/default
85ec6f318815a76cf8db71c0005d5f52  ./etc/hostname
8745e35fc02eb0c3830de728394681b2  ./usr/share/alsa/ucm/PandaBoardES/FMAnalog
876fc6cb1020e4655a366d7ac19c2b56  ./usr/share/alsa/pcm/surround40.conf
88bff9edee6b0465e03246541385b685  ./lib/libresolv-2.23.so
88d2957c29b72c341645dfff729ad078  ./usr/share/alsa/cards/ES1968.conf
897306b2f40bf7f1366aaf80001ddc8e  ./usr/share/alsa/init/00main
8b01179ed3318d963b0d08ccfad5ba56  ./usr/lib/libSDL2-2.0.so.0.4.1
8d182feb96489a2a6d1144ddeeccc42f  ./usr/lib/libtinyalsa.so.1.1.1
8faacf40a1e42e3089111d3cf64533dc  ./usr/share/alsa/cards/Maestro3.conf
8fac29e1233da15d1c1f3bf9379a159d  ./usr/lib/libSDL_sound-1.0.so.1.0.2
91768f15d2bdcd21972564048e4fe2e2  ./etc/shadow
92845205cbd87e5555d14815b3e6261b  ./lib/libnss_dns-2.23.so
939caf9384bf2d6922c6579bdbafa235  ./usr/share/alsa/ucm/skylake-rt286/Hdmi2
9528c60fae7c1680f86f2f20f9091192  ./etc/fstab
9531074fdcbf5fb106567f5113a08cdb  ./usr/share/alsa/cards/VIA8233.conf
95926cc32b581605465c7e3a5b515e9a  ./usr/share/alsa/cards/HDA-Intel.conf
95c7350bd921ce58b6f1b7dfbdb3ae1e  ./lib/libcrypt-2.23.so
96300e9c1b0dea5f61383f5d22342ef3  ./usr/share/terminfo/v/vt100
96e89c87b46837ad039a69a3825858b3  ./usr/share/alsa/ucm/DB410c/HDMI
99ec89843d154a619b6ce615a25b0603  ./usr/share/alsa/cards/Echo_Echo3G.conf
9a2f16ef24c59e19ef78243bec11cbfd  ./usr/share/alsa/cards/FWSpeakers.conf
9b2350fe94bbe3520766fc9af35f68fa  ./usr/lib/libavcodec.so.57.89.100
9ceb37badbc7870d320f2210b64b6242  ./usr/share/alsa/cards/Audigy2.conf
9eb8a2aa40e497b690b99a94982e0820  ./usr/share/alsa/topology/sklrt286/skl_i2s.conf
9ebdc6f90a6070a1cdac5ecde97fbd33  ./etc/services
a06a3fcbf8aebe420717a1933eb21572  ./usr/share/terminfo/v/vt100-putty
a0e8f0de39f34acf809a0b8110a65ade  ./usr/share/alsa/pcm/center_lfe.conf
a15e99b47a5596f2593399c71bf3bade  ./etc/protocols
a168b455fe0dd6b6b127c8441a011b45  ./lib/libnss_files-2.23.so
a1b52bd8139424509fb64a8e7968bd66  ./usr/share/alsa/cards/USB-Audio.conf
a2989220b87db6870b2933a64953e320  ./usr/lib/libjpeg.so.8.1.2
a44ccac429731eac5640030e423ad96b  ./usr/bin/playsound_simple
a81fae4ed8d1aa129e0412575bfa3ffa  ./usr/share/alsa/ucm/PAZ00/Record.conf
a854970c4f1e10a985affc5dcbfcb4ba  ./usr/share/alsa/cards/VXPocket.conf
a881925139729153a82420a80a6b763f  ./usr/share/udhcpc/default.script
a937461b77a1ab8dce43b14ddcde51f1  ./usr/bin/mplayer
aaedd7d9cb211efdc93f2f9028b9a667  ./usr/share/alsa/ucm/SDP4430/hifi
abf95b72cfd78d1f4ec9b5083105fba4  ./usr/share/alsa/cards/CMI8788.conf
ac77f8eb516b28d8fe9fb6896b83f1fb  ./usr/bin/tinymix
ac91545c5aee04892c03900dd933d4fd  ./usr/share/alsa/ucm/PandaBoard/hifiLP
adb71255b4b79c4b6dc3950b157c0f64  ./usr/lib/libavutil.so.55.58.100
ae4bfcc5e3c9031c501c6b695e8382d2  ./usr/lib/libSDL_ttf-2.0.so.0.10.1
ae6d7db349c2da7f0ddd38035e162166  ./usr/share/alsa/pcm/dpl.conf
affccadcedc40d4f0b98ec2e1077230c  ./usr/share/alsa/cards/ATIIXP.conf
b0ab22f448a751d9fcd5d75cb2eb5bae  ./usr/lib/libnl-idiag-3.so.200.25.0
b240423f9484099f7cb276d6b78b8fb8  ./usr/share/alsa/cards/RME9636.conf
b3ff28605cf956be4b3ca430ed37b979  ./usr/share/alsa/ucm/broxton-rt298/Hdmi2
b6dafb57b9fdebe9827e3c8d12dfa52f  ./usr/share/alsa/ucm/GoogleNyan/GoogleNyan.conf
b9581f156e0f3c8c2c1b45a674dc042f  ./usr/share/alsa/cards/ENS1371.conf
b9de29b647038e2f22d9b6671a8bcafb  ./usr/share/alsa/ucm/GoogleNyan/HiFi.conf
bab8daf32c0a604b08d8a06d52d0c001  ./usr/share/alsa/cards/Aureon51.conf
bed8278341f90f6124d1e8523a4a1a40  ./usr/share/alsa/ucm/PandaBoardES/voice
bf5091c9b3c6a5cb6d47a2b2363721e0  ./usr/lib/libavdevice.so.57.6.100
c2a58380d39779538b2f537b772f0583  ./usr/include/expat.h
c5d3f010ba554fb07c0466bbde04db10  ./usr/lib/libwebp.so.6.0.2
c645d95213274cb38d3ab1bf1ac054c5  ./usr/lib/libnl-route-3.so.200.25.0
ca3b114f0727da81a9b957b553a9915d  ./usr/share/terminfo/d/dumb
ca5a02ca1ef235995a59a3379e66f89b  ./usr/bin/playsound
cb1bd52cc2a0810ea01314040c1dae63  ./usr/share/alsa/pcm/dsnoop.conf
cbf8eddf89a7325b61575268044990c9  ./usr/share/alsa/pcm/front.conf
cc483958cd9471c830742ab261078b97  ./usr/share/alsa/ucm/PandaBoard/voice
ce023f2a7542adc65b1aadb71103bd24  ./usr/share/alsa/cards/VIA686A.conf
ce29e72b078978222f56828ed9f0a904  ./usr/share/alsa/cards/ATIIXP-SPDMA.conf
cfd25aeec5652fae1e7148eaf49666ce  ./usr/share/alsa/init/info
cfebdc2019b2b4b1d8f9a73a72d6630e  ./usr/lib/libSDL_gfx.so.13.9.1
d051a92017a865432db1f3fbdf1037e3  ./usr/share/alsa/ucm/PandaBoardES/record
d12535a6a1b7bde71e7e4d7b2c4b4f24  ./usr/share/alsa/cards/vc4-hdmi.conf
d1b30b4ff44177573e732fde3aba5a73  ./usr/lib/libz.so.1.2.11
d34d736572044528ba34505fc4b48ded  ./usr/share/alsa/pcm/surround71.conf
d3a9be1669e401e83e04cba82ef65f3e  ./usr/share/alsa/cards/SB-XFi.conf
d3be0607500338711912112f0cb3d690  ./usr/share/terminfo/x/xterm-color
d41d8cd98f00b204e9800998ecf8427e  ./dev/console
d49cef4e63470cb03a286d8fa7600902  ./usr/share/alsa/ucm/broxton-rt298/HiFi
d4a5a8b8a2691f0718df0cbb95ebc6d5  ./usr/lib/libnl-xfrm-3.so.200.25.0
d7950ce70c99d01c024288f6e78a7eb9  ./usr/share/alsa/cards/VXPocket440.conf
d7e7bcd59a397890b717536d863bd0e2  ./usr/share/alsa/ucm/PandaBoard/record
d7ff303130306098aed2bb874fff316f  ./lib/ld-2.23.so
d9aab39c652ff63dfca3a5828713cb34  ./usr/share/alsa/pcm/modem.conf
d9e3fe3df82dd7d9da384736f371e49e  ./usr/lib/libnl-3.so.200.25.0
dc866136f58d37621be73f1f0d8c8443  ./usr/share/alsa/cards/ICH.conf
dfb5dff22e7db68e5b7bc1b12fd8d155  ./usr/share/alsa/cards/AU8820.conf
e1b033dcbbd122cf88f2fd03ef403d55  ./usr/share/alsa/cards/CA0106.conf
e223aac5e591438a8cf631a0a6c9ddc7  ./lib/libnsl-2.23.so
e40b9ccb943a146281b5b44781e86b24  ./usr/lib/libexpat.la
e4111984063aec5feb2f0ed4ebe1a098  ./usr/share/alsa/cards/FireWave.conf
e4670570b66b6f803d3e88fd7ff7436b  ./usr/share/ffmpeg/libvpx-720p50_60.ffpreset
e4975c042c661ad8921996c9df2a6b84  ./usr/share/alsa/cards/TRID4DWAVENX.conf
e54a9cb37d75e2ec8a9a94131ee9a3e1  ./usr/share/ffmpeg/libvpx-360p.ffpreset
e59fb3dfb2ed8c2dbf0346cb949acfa0  ./usr/bin/aserver
e72221a58cd942ae77aa0a26287a0d99  ./usr/share/alsa/ucm/skylake-rt286/Hdmi1
eb2ad1b10bfab8ffcc5badacbacb2d1a  ./usr/lib/libform.so.6.0
eb560598005696f4fbb7047eb0947d99  ./usr/share/terminfo/s/screen-256color
ec8aa3e4ceb896858b84106fe8d7f10a  ./usr/lib/libtiff.so.5.2.6
ed1093b7d573caba3a233afc6dc47009  ./usr/share/alsa/cards/ICE1724.conf
ed1423f966d181095a576dd2cc0f53ae  ./usr/share/alsa/cards/AACI.conf
eeffa8ed7222fe2500deba90c92d5616  ./usr/lib/libSDL-1.2.so.0.11.4
efb6471ffe7fd0c8f873270fa9f0abc0  ./usr/bin/chocolate-doom
f0cefca49926d00903cf57551d901abe  ./usr/share/games/doom/doom1.wad
f1a15ab142b81e18d7b9698a11834940  ./usr/bin/alsamixer
f26497f940e6f4237a3b9727df750cb0  ./usr/share/alsa/cards/RME9652.conf
f2de602359d413373d5263d41c9bd87e  ./usr/lib/libavfilter.so.6.82.100
f2f7d9db5d9acd70f460c4cb86595014  ./usr/share/alsa/ucm/broxton-rt298/broxton-rt298.conf
f632f5cbce86412e9d973653b05e700d  ./etc/nsswitch.conf
f6fe06101e8056b64f6024e9ceb41cb3  ./usr/share/alsa/cards/CS46xx.conf
f755412d12bbc6970131fc9880dd5762  ./usr/share/alsa/cards/PC-Speaker.conf
f7d9d33439ea0d802f73ed227bb679b7  ./usr/lib/libpng16.so.16.31.0
f7f1886526b2461f6ed1cacba2a3ee01  ./usr/share/alsa/ucm/VEYRON-I2S/VEYRON-I2S.conf
fa515b54411bb81540c66a0d8e7fc479  ./usr/share/alsa/pcm/rear.conf
fb635260b61bcbddd954c3621317ccbd  ./etc/init.d/S01logging
fc1cc2be8a4062f0263c152c424a3684  ./usr/share/alsa/cards/CMI8338-SWIEC.conf
fc7b793918c16770379225d1451550b2  ./etc/hosts
fccb0ba4ee17416659246e88802b2f4d  ./usr/share/terminfo/s/screen
fd7609f2b678e68559d5a4fcd44f643b  ./usr/bin/aplay
fd7609f2b678e68559d5a4fcd44f643b  ./usr/bin/arecord
fdbc0e7424c5335e8cddde6d6244ed44  ./usr/bin/ffplay
fdd50693d25dbbc5fe5745325bb7032e  ./usr/lib/libavformat.so.57.71.100
fee3644c92735b3dbd736a63b1f1d8bd  ./usr/share/alsa/cards/CMI8338.conf
ffaa2685cd925369e0af65621a4716f1  ./lib/libdl-2.23.so
